# deployment/mydst-panel-update.service
[Unit]
Description=Update MyDST control panel from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash /opt/mydst/panel/deployment/panel-self-update.sh
TimeoutStartSec=45min
NoNewPrivileges=false

