# deployment/mydst-panel-update.path
[Unit]
Description=Watch for MyDST control panel update requests

[Path]
PathExists=/opt/mydst/panel-update.request
Unit=mydst-panel-update.service

[Install]
WantedBy=multi-user.target
